# systemd/master-ai-deep-clean.service
[Unit]
Description=Master AI biweekly deep clean (bug check + file clean)
Documentation=file:///home/elijah/scripts/deep_clean.sh
After=network-online.target

[Service]
Type=oneshot
ExecStart=/home/elijah/scripts/deep_clean.sh
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
# Don't murder the run if it takes a while — cleanup.sh can stretch
TimeoutStartSec=30min
